# deploy/gridbot.service
[Unit]
Description=GRID Trading Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=root
WorkingDirectory=/root/GRID_EXECUTE
EnvironmentFile=/root/GRID_EXECUTE/.env

# --resume so it picks up existing positions on restart
ExecStart=/root/GRID_EXECUTE/venv/bin/python3 live_trade.py --coin ETH --leverage 6 --resume

# Graceful shutdown: SIGTERM → runner.shutdown() → positions left open
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart on crash, but not on clean exit (e.g. circuit breaker halt)
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gridbot

[Install]
WantedBy=multi-user.target
